# container/base/s6-overlay/s6-rc.d/setup-persist/up
#!/command/execlineb -P
/usr/local/bin/setup-persist.sh
